# ffnh-led-control/files/etc/init.d/ffnh-led-control

#!/bin/sh /etc/rc.common
START=95
USE_PROCD=1
start_service(){ /usr/libexec/ffnh-led-control/apply.sh start; }
reload_service(){ /usr/libexec/ffnh-led-control/apply.sh apply; }
stop_service(){ /usr/libexec/ffnh-led-control/apply.sh stop; }
